# crontab
*/10 * * * * /bin/sh /home/user/Documents/dev/linux/scripts/expressVpnRandomizer.sh