# crontab
*/10 * * * * /bin/sh /replace/with/path/to/expressVpnRandomizer.sh
